# home/.config/systemd/user/albert.service
[Unit]
Description=Albert

[Service]
ExecStart=/usr/bin/albert

[Install]
WantedBy=user-applications.target
