# home/.config/systemd/user/compton.service
[Unit]
Description=xcompmgr

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=user-applications.target
