# home/.config/systemd/user/devmon.service
[Unit]
Description=Devmon

[Service]
ExecStart=/usr/bin/devmon

[Install]
WantedBy=user-applications.target
