# home/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/dropbox

[Install]
WantedBy=user-applications.target
