# home/.config/systemd/user/fcitx.service
[Unit]
Description=Flexible Context-aware Input Tool

[Service]
ExecStart=/usr/bin/fcitx5

[Install]
WantedBy=user-applications.target
