# home/.config/systemd/user/greenclip.service
[Unit]
Description=Greenclip

[Service]
ExecStart=%h/.local/bin/greenclip daemon

[Install]
WantedBy=user-applications.target
