# home/.config/systemd/user/light-locker.service
[Unit]
Description=Light Locker

[Service]
ExecStart=/usr/bin/light-locker

[Install]
WantedBy=user-applications.target
