# home/.config/systemd/user/nm-applet.service
[Unit]
Description=network monitor and control GUI applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=user-applications.target
