# home/.config/systemd/user/picom.service
[Unit]
Description=compositor

[Service]
ExecStart=/usr/bin/picom

[Install]
WantedBy=user-applications.target
