# home/.config/systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
ExecStart=/bin/polybar

[Install]
WantedBy=user-applications.target
