# home/.config/systemd/user/telegram.service
[Unit]
Description=Telegram Desktop

[Service]
ExecStart=/usr/bin/telegram-desktop -startintray

[Install]
WantedBy=user-applications.target
