# home/.config/systemd/user/trayer.service
[Unit]
Description=Trayer
After=xmobar.target

[Service]
ExecStart=%h/.local/bin/start-trayer

[Install]
WantedBy=xmobar.target
