# home/.config/systemd/user/trayer.service
[Unit]
Description=Trayer

[Service]
ExecStart=%h/.local/bin/start-trayer

[Install]
WantedBy=user-applications.target
