# home/.config/systemd/user/volumeicon.service
[Unit]
Description=VolumeIcon

[Service]
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=user-applications.target
