# home/.config/systemd/user/volumeicon.service
[Unit]
Description=VolumeIcon
After=pulseaudio.service

[Service]
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=user-applications.target
