# home/.config/systemd/user/xautolock.service
[Unit]
Description=Auto lock computer

[Service]
ExecStart=xautolock -locker 'dm-tool lock'

[Install]
WantedBy=user-applications.target
