# home/.config/systemd/user/xmobar.target
[Unit]
Description=Xmobar
Requires=user-applications.target
After=user-applications.target
