# home/.config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
ExecStart=/usr/bin/xremap %h/.config/xremap/config.yml

[Install]
WantedBy=user-applications.target
